# docker/cpu/conf.d/chromium.conf
[program:chromium]
command=chromium --no-sandbox file:///root/.local/share/jupyter/runtime/ > /dev/null
autorestart=true
