# docker/cpu/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared -rfbauth /app/x11vnc.pass 
autorestart=true
